# examples/migrate-from-s6/s6/s6-rc.d/api/finish
#!/command/execlineb -S0
foreground { echo "api exited with status ${1}" }
# s6 does not restart with backoff by default — you need to touch ./down
# and manage restarts manually, or rely on s6's default instant-restart loop.
exit 0
